# dockerfiles/supervisor/queue.conf
[program:queue]
command=/var/www/Laravel/artisan queue:work
numprocs=1
autostart=true
autorestart=true
