# dockerfiles/supervisor/websockets.conf
[program:websockets]
command=/var/www/Laravel/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
